# packaging/DEBIAN/prerm
#!/bin/sh
set -e

# Stop the service for the current user before removal.
systemctl --user stop clipd 2>/dev/null || true
systemctl --user disable clipd 2>/dev/null || true
