# systemd/clipd.service
[Unit]
Description=clipd - Clipboard History Manager
After=graphical-session.target ydotoold.service
Wants=ydotoold.service

[Service]
ExecStart=%h/.local/bin/clipd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
