# debian/postinst
#!/bin/sh
set -e

echo "clipd installed to /usr/local/bin/clipd"

if [ -n "$SUDO_USER" ]; then
    USER_UID=$(id -u "$SUDO_USER" 2>/dev/null) || USER_UID=""
    XDG_RT="/run/user/${USER_UID}"
    DBUS_SOCK="${XDG_RT}/bus"

    # Enable and start the systemd user services
    sudo -u "$SUDO_USER" \
        XDG_RUNTIME_DIR="$XDG_RT" \
        systemctl --user daemon-reload
    sudo -u "$SUDO_USER" \
        XDG_RUNTIME_DIR="$XDG_RT" \
        systemctl --user enable --now ydotoold || true
    sudo -u "$SUDO_USER" \
        XDG_RUNTIME_DIR="$XDG_RT" \
        systemctl --user enable --now clipd

    # Register the GNOME keyboard shortcut
    if [ -S "$DBUS_SOCK" ]; then
        sudo -u "$SUDO_USER" \
            XDG_RUNTIME_DIR="$XDG_RT" \
            DBUS_SESSION_BUS_ADDRESS="unix:path=${DBUS_SOCK}" \
            /usr/lib/clipd/setup-shortcut /usr/local/bin/clipd || true
    else
        # User not logged in — set the shortcut at next login via XDG autostart
        AUTOSTART_DIR="/home/${SUDO_USER}/.config/autostart"
        mkdir -p "$AUTOSTART_DIR"
        cp /etc/xdg/autostart/clipd-shortcut-setup.desktop "$AUTOSTART_DIR/" || true
        chown "$SUDO_USER" "$AUTOSTART_DIR/clipd-shortcut-setup.desktop" || true
    fi

    echo "clipd is running. Press Super+V to open the picker."
else
    echo ""
    echo "To start clipd for your user, run:"
    echo "  systemctl --user daemon-reload"
    echo "  systemctl --user enable --now clipd"
    echo ""
    echo "Then bind a key in your DE to: /usr/local/bin/clipd -trigger"
fi
