# debian/prerm
#!/bin/sh
set -e

if [ -n "$SUDO_USER" ]; then
    USER_UID=$(id -u "$SUDO_USER" 2>/dev/null) || USER_UID=""
    XDG_RT="/run/user/${USER_UID}"
    DBUS_SOCK="${XDG_RT}/bus"

    # Stop and disable the systemd user services
    sudo -u "$SUDO_USER" \
        XDG_RUNTIME_DIR="$XDG_RT" \
        systemctl --user disable --now clipd 2>/dev/null || true
    sudo -u "$SUDO_USER" \
        XDG_RUNTIME_DIR="$XDG_RT" \
        systemctl --user disable --now ydotoold 2>/dev/null || true

    # Remove the GNOME keyboard shortcut
    if [ -S "$DBUS_SOCK" ]; then
        sudo -u "$SUDO_USER" \
            XDG_RUNTIME_DIR="$XDG_RT" \
            DBUS_SESSION_BUS_ADDRESS="unix:path=${DBUS_SOCK}" \
            /usr/lib/clipd/teardown-shortcut || true
    fi

    # Remove the per-user autostart fallback if it was placed there
    rm -f "/home/${SUDO_USER}/.config/autostart/clipd-shortcut-setup.desktop" || true
else
    systemctl --user stop clipd 2>/dev/null || true
    systemctl --user disable clipd 2>/dev/null || true
fi
